# etc/init.d/ip6neigh
#!/bin/sh /etc/rc.common

START=96

readonly SBIN_PATH="/usr/sbin/"
readonly SERVICE_NAME="ip6neigh_svc.sh"
readonly SERVICE_SCRIPT="${SBIN_PATH}${SERVICE_NAME}"

errormsg() {
	local msg="Error: $1"
	>&2 echo "$msg"
	logger -t ip6neigh "$msg"
	exit 1
}

start() {
	pgrep -f "$SERVICE_NAME" >/dev/null && return 0
	
	[ -f "$SERVICE_SCRIPT" ] || errormsg "The main script file $SERVICE_SCRIPT is missing."

	SERVICE_DAEMONIZE=1
	service_start "$SERVICE_SCRIPT" -s

	return 0
}

stop() {
	killall "$SERVICE_NAME" 2>/dev/null
	return 0
}

restart() {
	stop
	sleep 1
	start
	return 0
}
