# etc/init.d/ip6neigh
#!/bin/sh /etc/rc.common

#ip6neigh init.d script

START=96

readonly VERSION="2.0.0"
readonly SBIN_PATH="/usr/sbin/"
readonly SERVICE_NAME="ip6neigh-svc.sh"
readonly SERVICE_SCRIPT="${SBIN_PATH}${SERVICE_NAME}"
readonly CONFIG_FILE="/etc/config/ip6neigh"

. /lib/functions.sh

#Loads UCI configuration file
[ -f "$CONFIG_FILE" ] || errormsg "The UCI config file $CONFIG_FILE is missing. A template is avaliable at https://github.com/AndreBL/ip6neigh ."
reset_cb
config_load ip6neigh
config_get LAN_IFACES config lan_iface

errormsg() {
	local msg="Error: $1"
	>&2 echo "$msg"
	logger -t ip6neigh "$msg"
	exit 1
}

start() {
	pgrep -f "$SERVICE_NAME" >/dev/null && return 0
	
	[ -f "$SERVICE_SCRIPT" ] || errormsg "The main script file $SERVICE_SCRIPT is missing."

	local snoop=$(uci get ip6neigh.config.dad_snoop 2>/dev/null)

	SERVICE_DAEMONIZE=1

	for LAN_IFACE in $LAN_IFACES ; do
		service_start "$SERVICE_SCRIPT" --interface="$LAN_IFACE"
	done

	[ "$snoop" = '1' ] && service_start "$SERVICE_SCRIPT" -n

	return 0
}

stop() {
	local ppid
	local cpid
	
	pgrep -f "$SERVICE_NAME" |
		while IFS= read -r ppid
		do
			pgrep -P "$ppid" |
				while IFS= read -r cpid
				do
					kill "$cpid" 2>/dev/null
				done
			kill "$ppid" 2>/dev/null
		done

	rm -f /tmp/ip6neigh.*.cache
	rm -f /tmp/hosts/ip6neigh.*
	killall -1 dnsmasq

	return 0
}

restart() {
	stop
	sleep 1
	start
	return 0
}
