# etc/init.d/ip6neigh
#!/bin/sh /etc/rc.common

START=96

readonly SCRIPT_FILE="/root/ip6neigh_mon.sh"

start() {
	pgrep -f ip6neigh_mon.sh && return 0

	SERVICE_DAEMONIZE=1
	service_start "$SCRIPT_FILE" -s

	return 0
}

stop() {
	killall ip6neigh_mon.sh 2>/dev/null
	return 0
}

restart() {
	stop
	sleep 1
	start
	return 0
}
